# meta-ros2-autom8/recipes-support/startup-script/files/startup-script-001.service
[Unit]
Description=A test systemd service

[Service]
ExecStart=/bin/echo "Hello"

[Install]
WantedBy=multi-user.target
